# pi/church-translator.service
[Unit]
Description=Church Translator Relay Server
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/church-translator-relay /etc/church-translator/config.yaml
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
